# assets/performanceprofile/configs/dedicatedcpus.slice
[Unit]
Description=Top level slice for dedicated CPUs used by services/applications that requires full isolation.

[Slice]
AllowedCPUs={{ .DedicatedCpus }}
