# assets/performanceprofile/configs/ovs.slice
[Unit]
Description=Top level slice used to give openvswitch access to an unrestricted set of cpus

[Slice]
